# Code/out/_lib_systemd_system_/outCam.service
[Unit]
Description=Sanders Outdoor Conditions Cam  camStream.service
After=network-online.target

[Service]
Type=simple
ExecStart=/home/greg/out/camStream.py
WorkingDirectory=/home/greg/out/
StandardOutput=syslog
StandardError=syslog
#User=greg
#Group=www-data
#ExecStop = /bin/kill -2 $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-abort

[Install]
WantedBy=multi-user.target
