# Code/out/_lib_systemd_system_/outSocket.service
[Unit]
Description=Sanders WeatherPi Unix Socket Utility
After=network-online.target

[Service]
Type=simple
ExecStart=/home/greg/out/outSocket.py
WorkingDirectory=/home/greg/out/
StandardOutput=syslog
StandardError=syslog
User=root
Group=root
#ExecStop = /bin/kill -2 $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-abort

[Install]
WantedBy=multi-user.target
