# Code/out/_lib_systemd_system_/rainMainDATA.service
[Unit]
Description=Sanders WeatherPi Rainfall Data Logger  rainMainDATA.service
After=network-online.target

[Service]
Type=simple
ExecStart=/home/YourUserName/out/rainMainDATA.py
WorkingDirectory=/home/YourUserName/out/
StandardOutput=syslog
StandardError=syslog
User=YourUserName
Group=aGroupYouAreMemberOf(like'pi')
#Group=pi
#ExecStop = /bin/kill -2 $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-abort

[Install]
WantedBy=multi-user.target
